# git.service
[Unit]
Description=Git-Web Service
After=syslog.target

[Service]
Type=simple
ExecStart=/etc/init.d/git

[Install]
WantedBy=multi-user.target
